# webtop.service
[Unit]
Description=webtop - lightweight web-based system dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/webtop -port 9999
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
